# src/services/picoclaw-web.service
[Unit]
Description=Picoclaw Web UI (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=stas
WorkingDirectory=/home/stas/.picoclaw
ExecStart=/usr/bin/python3 /home/stas/.picoclaw/bot_web.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/stas/.picoclaw/bot.env
Environment=PICOCLAW_BIN=/usr/bin/picoclaw

[Install]
WantedBy=multi-user.target
